# debian/quantum-server.postrm
#!/bin/sh -e

if [ "$1" = purge ]; then
	update-rc.d quantum-server remove >/dev/null
fi

#DEBHELPER#
