# debian/quantum-dhcp-agent.postrm
#!/bin/sh -e

if [ "$1" = purge ]; then
	update-rc.d quantum-plugin-dhcp remove >/dev/null
fi

#DEBHELPER#
